# supervisord.conf
[unix_http_server]
file=/run/supervisor/supervisor.sock

[supervisord]
logfile=/var/log/supervisor/supervisord.log
logfile_maxbytes=10MB
pidfile=/run/supervisor/supervisord.pid
nodaemon=true
user=root

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[supervisorctl]
serverurl=unix:///run/supervisor/supervisor.sock

# websockify — token routing for VNC and audio
[program:websockify]
command=/usr/bin/websockify --web=/usr/share/novnc --token-plugin=TokenFile --token-source=/etc/websockify/token.cfg 0.0.0.0:8080
autorestart=true
priority=10
stdout_logfile=/var/log/supervisor/websockify.log
stderr_logfile=/var/log/supervisor/websockify.log

# audio-proxy — GStreamer pipeline encoding PulseAudio to WebM/Opus
[program:audio-proxy]
command=/usr/local/bin/audio-proxy.sh -l 5711
autorestart=true
priority=15
startsecs=5
stdout_logfile=/var/log/supervisor/audio-proxy.log
stderr_logfile=/var/log/supervisor/audio-proxy.log

# QEMU
[program:qemu]
command=/usr/local/bin/start-qemu.sh
autorestart=false
priority=20
startsecs=3
stdout_logfile=/var/log/supervisor/qemu.log
stderr_logfile=/var/log/supervisor/qemu.log

# setup-audio — loads PulseAudio TCP module in museum VM via SSH
[program:setup-audio]
command=/usr/local/bin/setup-audio.sh
autorestart=false
priority=45
startsecs=0
stdout_logfile=/var/log/supervisor/setup-audio.log
stderr_logfile=/var/log/supervisor/setup-audio.log
