# supervisord.conf
[unix_http_server]
file=/run/supervisor/supervisor.sock

[supervisord]
logfile=/var/log/supervisor/supervisord.log
logfile_maxbytes=10MB
pidfile=/run/supervisor/supervisord.pid
nodaemon=true
user=root

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[supervisorctl]
serverurl=unix:///run/supervisor/supervisor.sock

# PulseAudio — null sink for QEMU audio output
[program:pulseaudio]
command=/usr/bin/pulseaudio --exit-idle-time=-1 --log-target=stderr --log-level=0
autorestart=true
priority=5
environment=HOME="/root",XDG_RUNTIME_DIR="/tmp/pulse",PULSE_RUNTIME_PATH="/tmp/pulse"
stdout_logfile=/var/log/supervisor/pulseaudio.log
stderr_logfile=/var/log/supervisor/pulseaudio.log

# websockify — token routing for VNC and audio
[program:websockify]
command=/usr/bin/websockify --web=/usr/share/novnc --token-plugin=TokenFile --token-source=/etc/websockify/token.cfg 0.0.0.0:8080
autorestart=true
priority=10
stdout_logfile=/var/log/supervisor/websockify.log
stderr_logfile=/var/log/supervisor/websockify.log

# audio-proxy — GStreamer pipeline encoding PulseAudio to WebM/Opus
[program:audio-proxy]
command=/usr/local/bin/audio-proxy.sh -l 5711
autorestart=true
priority=15
startsecs=5
environment=PULSE_SERVER="unix:/tmp/pulse/native",HOME="/root",PULSE_RUNTIME_PATH="/tmp/pulse"
stdout_logfile=/var/log/supervisor/audio-proxy.log
stderr_logfile=/var/log/supervisor/audio-proxy.log

# QEMU
[program:qemu]
command=/usr/local/bin/start-qemu.sh
autorestart=false
priority=20
startsecs=5
environment=PULSE_SERVER="unix:/tmp/pulse/native",HOME="/root"
stdout_logfile=/var/log/supervisor/qemu.log
stderr_logfile=/var/log/supervisor/qemu.log
