# supervisord.conf
[unix_http_server]
file=/run/supervisor/supervisor.sock

[supervisord]
logfile=/var/log/supervisor/supervisord.log
logfile_maxbytes=10MB
pidfile=/run/supervisor/supervisord.pid
nodaemon=true
user=root

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[supervisorctl]
serverurl=unix:///run/supervisor/supervisor.sock

# noVNC proxies QEMU's built-in VNC (port 5901) to browser WebSocket (port 8080)
[program:novnc]
command=/usr/share/novnc/utils/novnc_proxy --vnc 127.0.0.1:5901 --listen 0.0.0.0:8080
autorestart=true
priority=10
stdout_logfile=/var/log/supervisor/novnc.log
stderr_logfile=/var/log/supervisor/novnc.log

# QEMU — built-in VNC on display :1 (port 5901)
[program:qemu]
command=/usr/local/bin/start-qemu.sh
autorestart=false
priority=20
startsecs=3
stdout_logfile=/var/log/supervisor/qemu.log
stderr_logfile=/var/log/supervisor/qemu.log
